# command/etc/profile.d/disable-builtins.sh
enable -n help
enable -n echo
enable -n jobs
enable -n kill
enable -n logout
enable -n set
enable -n suspend
enable -n type
enable -n umask
enable -n wait
enable -n unalias
enable -n bind
enable -n builtin
enable -n command
enable -n declare
enable -n dirs
enable -n disown

# 必ず最後に実行する
enable -n enable